# roles/ubuntu_clients/configure_lehrerkiosk/configure_desktop-connect2home/files/set-icon-positions_connect2home.conf
start on desktop-start
stop on desktop-end

script
gvfs-set-attribute -t string 'Schreibtisch/connect2home.desktop' 'metadata::nautilus-icon-position' '1400,100'
end script

