# roles/ubuntu_clients/configure_lehrerkiosk/configure_desktop-serverconnections-rw_in_userspace/files/set-icon-positions-connect2video-rw.conf
start on desktop-start
stop on desktop-end

script
gvfs-set-attribute -t string 'Schreibtisch/connect2video-rw.desktop' 'metadata::nautilus-icon-position' '1400,200'
end script

