# roles/ubuntu_clients_20.04/configure_lehrerkiosk/configure_launcher-webuntis4localuser/templates/set-apps-in-launcher.conf.j2
start on desktop-start
stop on desktop-end

script
gsettings set com.canonical.Unity.Launcher favorites "['application://ubiquity.desktop', 'application://webuntis.desktop', 'application://org.gnome.Nautilus.desktop', 'application://firefox.desktop', 'application://geogebra.desktop', 'application://vlc.desktop', 'application://libreoffice-writer.desktop', 'application://libreoffice-impress.desktop', 'application://unity-control-center.desktop', 'unity://running-apps', 'unity://desktop-icon', 'unity://expo-icon', 'unity://devices']"
end script

